# supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=debug

[program:nginx]
command=nginx -g "daemon off;"
user=root

[program:uvicorn]
directory=/app/backend/
command=uvicorn main:app --host 0.0.0.0 --port 8000
user=app-import-tool
stdout_logfile=/dev/stdout
redirect_stderr=true
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0